# {{ cookiecutter.project_name }}/deploy/etc/s6-overlay/s6-rc.d/worker/run
#!/command/with-contenv execlineb

{{ cookiecutter.project_name }} qcluster