# {{ cookiecutter.project_name }}/deploy/etc/s6-overlay/s6-rc.d/server/run
#!/command/with-contenv execlineb

cd /app
exec python {{ cookiecutter.project_name }}/__main__.py litestream replicate -exec "python {{ cookiecutter.project_name }}/__main__.py"
